# systemd/protonlumoai.service
[Unit]
Description=ProtonLumoAI - Processeur d'Emails Intelligent
After=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/ProtonLumoAI
ExecStart=/home/%i/ProtonLumoAI/run.fish
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
