# systemd/protonlumoai.service
[Unit]
Description=ProtonLumoAI - Automated Email Classifier with Adaptive Learning
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/ProtonLumoAI/scripts
Environment="PATH=%h/ProtonLumoAI/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=%h/ProtonLumoAI/scripts"
EnvironmentFile=%h/ProtonLumoAI/.env
ExecStartPre=/bin/sleep 15
ExecStart=%h/ProtonLumoAI/venv/bin/python %h/ProtonLumoAI/scripts/email_processor.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=protonlumoai

# Protection et sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/ProtonLumoAI/data %h/ProtonLumoAI/logs %h/ProtonLumoAI/config

[Install]
WantedBy=default.target